# systemd/tcdd.service
[Unit]
Description=TCDD Detection System
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=tcdd-thesis
WorkingDirectory=/home/tcdd-thesis/tcdd-dev
ExecStartPre=/bin/chmod +x /home/tcdd-thesis/tcdd-dev/start.sh
ExecStart=/bin/bash /home/tcdd-thesis/tcdd-dev/start.sh
Environment=DISPLAY=:0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target